# src/cron.d/01-apps-api-addpackages
*/1 * * * * repo /opt/apps-api/background_jobs/addpackages.sh
 