# src/cron.d/01-apps-api-bgtasks
30 6 * * 1 root systemctl start apps-api-bgtask-github.service
30 6 * * * root systemctl start apps-api-bgtask-duolingo.service