# src/config/cron.d/01-apps-api-jsonstats
*/1 * * * * apps /opt/apps-api/bin/jsonstats
